# internal/testdata/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>{{.BundleID}}</string>
        <key>EnvironmentVariables</key>
        <dict>
            <key>GO_TEST_SERVER_ADDR</key>
            <string>{{.GoTestServerAddr}}</string>
        </dict>
        <!--Run-Test-Binary-->
        <key>Program</key>
        <string>{{.GoTestBinary}}</string>
        <key>ProgramArguments</key>
        <array>
            <!--Run-Only-Remote-Tests -->
            <string>-test.count=1</string>
            <string>-test.run={{.GoTestName}}</string>
            <string>-test.timeout=30s</string>
            <string>-test.v=true</string>
{{if .GoCoverDir}}
            <string>-test.gocoverdir={{.GoCoverDir}}</string>
{{else}}
            <!--cover-flag-not-specified-->
{{end}}
        </array>
        <key>RunAtLoad</key>
        <true />
        <key>StandardErrorPath</key>
        <string>{{.StderrFile}}</string>
        <key>StandardOutPath</key>
        <string>{{.StdoutFile}}</string>
        <!--Sockets-->
        <key>Sockets</key>
        <dict>
            <!--tcp-->
            <key>tcp</key>
            <dict>
                <key>SockServiceName</key>
                <string>{{.TCP}}</string>
                <key>SockType</key>
                <string>stream</string>
                <key>SockFamily</key>
                <string>IPv4</string>
                <key>SockNodeName</key>
                <string>localhost</string>
            </dict>
            <!--udp-->
            <key>udp</key>
            <dict>
                <key>SockServiceName</key>
                <string>{{.UDP}}</string>
                <key>SockType</key>
                <string>dgram</string>
                <key>SockFamily</key>
                <string>IPv4</string>
                <key>SockNodeName</key>
                <string>localhost</string>
            </dict>
            <!--tcp-multiple-->
            <key>tcp-multiple</key>
            <dict>
                <key>SockServiceName</key>
                <string>{{.TCPMultiple}}</string>
                <key>SockType</key>
                <string>stream</string>
                <key>SockNodeName</key>
                <string>localhost</string>
            </dict>
            <!--udp-multiple-->
            <key>udp-multiple</key>
            <dict>
                <key>SockServiceName</key>
                <string>{{.UDPMultiple}}</string>
                <key>SockType</key>
                <string>dgram</string>
                <key>SockNodeName</key>
                <string>localhost</string>
            </dict>
            <!--tcp-dualstack-single-socket-->
            <key>tcp-dualstack-single-socket</key>
            <dict>
                <key>SockServiceName</key>
                <string>{{.TCPDualStackSingleSocket}}</string>
                <key>SockType</key>
                <string>stream</string>
                <key>SockFamily</key>
                <string>IPv4v6</string>
                <key>SockNodeName</key>
                <string>localhost</string>
            </dict>
            <!--udp-dualstack-single-socket-->
            <key>udp-dualstack-single-socket</key>
            <dict>
                <key>SockServiceName</key>
                <string>{{.UDPDualStackSingleSocket}}</string>
                <key>SockType</key>
                <string>dgram</string>
                <key>SockFamily</key>
                <string>IPv4v6</string>
                <key>SockNodeName</key>
                <string>localhost</string>
            </dict>
            <!--tcp-invalid-type-->
            <key>tcp-invalid-type</key>
            <dict>
                <key>SockServiceName</key>
                <string>{{.TCPInvalidType}}</string>
                <key>SockType</key>
                <string>dgram</string>
                <key>SockFamily</key>
                <string>IPv4</string>
                <key>SockNodeName</key>
                <string>localhost</string>
            </dict>
            <!--udp-invalid-type-->
            <key>udp-invalid-type</key>
            <dict>
                <key>SockServiceName</key>
                <string>{{.UDPInvalidType}}</string>
                <key>SockType</key>
                <string>stream</string>
                <key>SockFamily</key>
                <string>IPv4</string>
                <key>SockNodeName</key>
                <string>localhost</string>
            </dict>
            <!--unix-stream-socket-->
            <key>unix-stream</key>
            <dict>
                <key>SockPathName</key>
                <string>{{.UnixStreamSocket}}</string>
                <key>SockPathMode</key>
                <integer>448</integer>
                <key>SockType</key>
                <string>stream</string>
            </dict>
            <!--unix-datagram-socket-->
            <key>unix-datagram</key>
            <dict>
                <key>SockPathName</key>
                <string>{{.UnixDatagramSocket}}</string>
                <key>SockPathMode</key>
                <integer>448</integer>
                <key>SockType</key>
                <string>dgram</string>
            </dict>
        </dict>
    </dict>
</plist>
